# deploy/supervisord.conf
[program:djblog]
command=pipenv run gunicorn djblog.wsgi
directory=/data/deploy_app/djblog
user=deploy
autostart=true
autorestart=true
redirect_stderr=true
